# config/launchd/com.mcrouch.ntfy-claude-code.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcrouch.ntfy-claude-code</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/.config/agent-notify/ntfy-subscribe.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mcrouch/Library/Logs/ntfy-claude-code.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mcrouch/Library/Logs/ntfy-claude-code.log</string>
</dict>
</plist>
